# whonix_shared/etc/profile.d/20_mesg.sh
#!/bin/bash

## Whonix /etc/profile.d/20_mesg.sh

## Gets run with any login shell.

#set -x

#mesg

## Allow messages to tty
mesg y

#mesg

## End of Whonix /etc/profile.d/20_mesg.sh
