# whonix_shared/etc/profile.d/30_timesync.sh
#!/bin/bash

## Whonix /etc/profile.d/30_timesync.sh

## Gets run with any login shell.

/usr/lib/whonix/delay timesync --autostart &

## End of Whonix /etc/profile.d/30_timesync.sh
