# whonix_shared/etc/profile.d/40_whonixcheck.sh
#!/bin/bash

## Whonix /etc/profile.d/40_whonixcheck.sh

## Gets run with any login shell.

/usr/lib/whonix/delay whonixcheck --autostart &

## End of Whonix /etc/profile.d/40_whonixcheck.sh
