# whonix_workstation/etc/profile.d/20_whonix.sh
#!/bin/bash

## TorBirdy reads this variable and sets appropriate socks proxy settings for
## Stream Isolation.
export WHONIX=1 

