# whonix_shared/etc/cron.hourly/40_whonixcheck
#!/bin/bash

## Whonix /etc/cron.hourly/40_whonixcheck

set -x

## Debugging:
#echo "whonixcheck cron.hourly pre: date: $(date)" >> /home/user/delayed

## & did not work on Debian.
exec /usr/bin/whonixcheck_hourly

## Debugging:
#echo "whonixcheck cron.hourly post: date: $(date)" >> /home/user/delayed

exit 0

## End of Whonix /etc/cron.hourly/40_whonixcheck
