# whonix_shared/etc/cron.hourly/htpdate
#!/bin/bash

## Whonix /etc/cron.hourly/htpdate

set -x

## Debugging:
#echo "htpdate cron.hourly pre: date: `date`" >> /home/user/htpdatelog

## & did not work on Debian.
exec /usr/local/bin/htpdate_hourly

## Debugging:
#echo "htpdate cron.hourly post: date: `date`" >> /home/user/htpdatelog

exit 0

## End of Whonix /etc/cron.hourly/htpdate

