# whonix_shared/etc/cron.hourly/whonixcheck
#!/bin/bash

## Whonix /etc/cron.hourly/whonixcheck

set -x

## Debugging:
#echo "whonixcheck cron.hourly: date: `date`" >> /home/user/delayed

## & did not work on Debian.
exec /usr/local/bin/whonixcheck_hourly

exit 0

## End of Whonix /etc/cron.hourly/whonixcheck
