# aos_workstation/home/user/.config/autostart/aoscheck.desktop
[Desktop Entry]
Type=Application
Exec=aoscheck 
