# aos_workstation/home/user/.config/autostart/torcheck.desktop
[Desktop Entry]
Type=Application
Exec=torcheck 
