# whonix_workstation/home/user/.config/autostart/whonixcheck.desktop
[Desktop Entry]
Type=Application
Exec=whonixcheck_login
