# whonix_shared/etc/xdg/autostart/timesync.desktop
[Desktop Entry]
Type=Application
Exec=timesync --autostart
