# whonix_shared/etc/xdg/autostart/whonixcheck.desktop
[Desktop Entry]
Type=Application
Exec=whonixcheck --showx
